# etc/services.d/ssl_sync/run
#!/usr/bin/with-contenv bash
## file: etc/services.d/ssl_sync/run
# Этот скрипт запускается s6-overlay как longrun-сервис

# Перенаправляем весь вывод в логи аддона
exec > /proc/1/fd/1 2>/proc/1/fd/2

echo "[$(date '+%Y-%m-%d %H:%M:%S')] SSL Sync сервис запущен — начинаем бесконечный цикл"

while true; do
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] Запуск синхронизации сертификатов…"
    
    # ← ВОТ ГЛАВНОЕ ИСПРАВЛЕНИЕ: используем exec!
    exec /run.sh
    
    # Если по какой-то причине exec не сработает — логируем и спим
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] Ошибка выполнения run.sh, ждём 30 секунд перед повтором"
    sleep 30
done
