# etc/services.d/ssl_sync/run
#!/usr/bin/with-contenv bash
## file: etc/services.d/ssl_sync/run
# Этот скрипт запускается s6-overlay и работает постоянно

# Перенаправляем stdout/stderr в логи аддона
exec > /proc/1/fd/1 2>/proc/1/fd/2

echo "[$(date '+%Y-%m-%d %H:%M:%S')] Сервис SSL Sync запущен и работает в фоне"

# Бесконечный цикл с периодическим запуском синхронизации
while true; do
    # Запускаем основной скрипт синхронизации
    /run.sh

    # Пауза между запусками (в секундах)
    # 300 = каждые 5 минут, можно поставить 3600 = каждый час
    sleep 300
done
