# etc/services.d/ssl_sync/run
#!/bin/sh
# Скрипт запускается как сервис PID 1 через s6
/run.sh
