# etc/services.d/ssl_sync/run
#!/usr/bin/with-contenv bash
# Этот скрипт запускает ваш сервис под s6-overlay

# Перенаправляем весь вывод в логи аддона
exec 2>&1

echo "[$(date '+%Y-%m-%d %H:%M:%S')] SSL Sync сервис запущен — начинаем цикл каждые 5 минут"

while true; do
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] → Запуск синхронизации сертификатов"
    
    # ← ЭТО ОБЯЗАТЕЛЬНО! exec заменяет текущий процесс на run.sh
    exec /run.sh
    
    # Если exec по какой-то причине провалится (никогда не должно), спим и пробуем снова
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] Ошибка при выполнении run.sh — повтор через 60 сек"
    sleep 60
done

fix(s6): use exec /run.sh to avoid suexec error
