# etc/services.d/ssl_sync/run
#!/usr/bin/with-contenv bash
exec > /proc/1/fd/1 2>/proc/1/fd/2
echo "[$(date '+%Y-%m-%d %H:%M:%S')] SSL Sync service started"

while true; do
    /run.sh
    sleep 300
done

